# ifc-bonsai-mcp-main/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

; ─────────────────────────────────────────────────────────
; Process 1: Blender running headless with the MCP addon
; Must start FIRST so the socket on port 9876 is ready
; before the MCP SSE server tries to connect to it.
; ─────────────────────────────────────────────────────────
[program:blender]
command=blender -b
    --addons blendermcp
    --python /app/docker/blender_autostart.py
priority=10
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/blender.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/blender_err.log
stderr_logfile_maxbytes=5MB
environment=DISPLAY="",HOME="/root",BLENDER_USER_CONFIG="/app/.blender_config"

; ─────────────────────────────────────────────────────────
; Process 2: Embedding server for RAG
; Runs on port 9090 using the pre-built offline cache
; ─────────────────────────────────────────────────────────
[program:embedding-server]
command=python scripts/embedding_server.py
    --model sentence-transformers/all-MiniLM-L6-v2
    --host 127.0.0.1
    --port 9090
    --normalize
    --offline
    --cache-dir /app/.cache/huggingface
priority=20
autostart=true
autorestart=true
startretries=5
startsecs=10
directory=/app
stdout_logfile=/var/log/embedding_server.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/embedding_server_err.log
stderr_logfile_maxbytes=5MB
environment=PYTHONPATH="/app/src",PYTHONUNBUFFERED="1",HF_HUB_OFFLINE="1",TRANSFORMERS_OFFLINE="1"

; ─────────────────────────────────────────────────────────
; Process 3: MCP SSE HTTP server
; Starts AFTER Blender (priority=30) and is what App Runner
; health-checks and routes traffic to on port 8080
; ─────────────────────────────────────────────────────────
[program:mcp-server]
command=python /app/docker/serve_sse.py
priority=30
autostart=true
autorestart=true
startretries=5
startsecs=15
directory=/app
stdout_logfile=/var/log/mcp_server.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/mcp_server_err.log
stderr_logfile_maxbytes=5MB
environment=PYTHONPATH="/app/src",PYTHONUNBUFFERED="1",BLENDER_MCP_HOST="127.0.0.1",BLENDER_MCP_PORT="9876",BLENDER_MCP_REMOTE_EMBEDDINGS_URL="http://127.0.0.1:9090/embeddings",BLENDER_MCP_EMBEDDING_OFFLINE="1",HF_HUB_OFFLINE="1"
