# ifc-bonsai-mcp-main/docker/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/dev/fd/1
logfile_maxbytes=0
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:blender]
command=/usr/local/bin/blender -b --addons blendermcp --python /app/docker/blender_autostart.py
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:embedding_server]
command=python scripts/embedding_server.py --model sentence-transformers/all-MiniLM-L6-v2 --cache-dir /app/.cache/huggingface --host 127.0.0.1 --port 9090 --normalize --offline
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:mcp_sse_server]
command=python /app/docker/serve_sse.py
autostart=true
autorestart=true
priority=30
environment=PYTHONPATH="/app/src",PORT="8000"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
